# debian/turbod_0.1.0-1_amd64/etc/systemd/system/turbod.service
[Install]
WantedBy = multi-user.target

[Unit]
Description = Turbo daemon

[Service]
Type = simple
After = network-online.target
ExecStart = /usr/bin/turbod
Restart=on-failure